# Server Rebuild/DeathStar/Systemd Timer/FileSync.service
[Unit]
Description=FileSync

[Service]
Type=oneshot
