# Server Rebuild/DeathStar/Systemd Timer/video_sync.timer
[Unit]
Description=video_sync

[Timer]
OnCalendar=00/6:50
Persistent=true

[Install]
WantedBy=timers.target