# Archive/Digital Ocean Rebuild/DeathStar/Systemd Timer/Podcast_Sync.timer
[Unit]
Description=Podcast_Sync

[Timer]
OnCalendar=00/1:50
Persistent=true

[Install]
WantedBy=timers.target