# Server Rebuild/Fulminatrix/Systemd Timer/Podman_Cleanup.service
[Unit]
Description=Podman_Cleanup

[Service]
Type=oneshot
ExecStart=podman rm --all
WorkingDirectory=/